# Layout/DEBIAN/postinst
#!/bin/sh
echo "Loading ProxySwitcher daemon..."
launchctl load /Library/LaunchDaemons/com.esehr.proxyswitcherd.plist
exit 0